# docker/app/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

; ─────────────────────────────────────────────────────────────────────────────
; Laravel Octane (RoadRunner) – main HTTP server
; ─────────────────────────────────────────────────────────────────────────────
[program:octane]
command=php /var/www/html/artisan octane:start
    --server=roadrunner
    --host=0.0.0.0
    --port=8000
    --workers=4
    --max-requests=500
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─────────────────────────────────────────────────────────────────────────────
; Queue workers (2 parallel workers)
; ─────────────────────────────────────────────────────────────────────────────
[program:queue]
command=php /var/www/html/artisan queue:work
    --sleep=3
    --tries=3
    --max-time=3600
    --queue=default
directory=/var/www/html
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
